# docker/root/etc/services.d/messenger-consumer/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eux

cd /app

exec 2>&1
exec su-exec www-data:www-data bin/console messenger:consume -vv async
