# docker/root/etc/services.d/messenger-consumer/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eux

cd /app

unset HOME

exec 2>&1
exec s6-setuidgid www-data bin/console messenger:consume -vv async
