# docker/root/etc/services.d/nginx-error-log/run
#!/usr/bin/env sh

set -eux

exec logutil-service -f /run/s6/nginx-error-log.fifo /var/log/nginx-error-log
