# docker/root/etc/services.d/caddy/run
#!/command/with-contenv sh
# shellcheck shell=sh

exec caddy run --config config.yaml --adapter yaml
