# docker/root/etc/services.d/workerman/finish
#!/usr/bin/env sh

set -eux

exec /run/s6/basedir/bin/halt
