# docker/root/etc/services.d/ofelia/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eux

exec su-exec www-data:www-data ofelia daemon --config=config.ini
