# docker/root/etc/services.d/ofelia/run
#!/command/with-contenv sh
# shellcheck shell=sh

set -eux

exec s6-setuidgid www-data ofelia daemon --config=config.ini
